# automation/stock-fetch.service
[Unit]
Description=Stock Derived Data Maintenance Service
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/var/www/valuation
Environment="PATH=/var/www/valuation/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
# Load valuation and cache configuration from project .env
EnvironmentFile=/var/www/valuation/.env
# Refresh security master and derived valuation cache only.
ExecStart=/var/www/valuation/.venv/bin/python3 /var/www/valuation/run_pipeline.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
